# platform/deploy/systemd/jobsvc.service
# jobsvc.service — FastAPI job service.
#
# Install:
#   sudo install -m 644 jobsvc.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now jobsvc.service
#
# Logs:
#   sudo journalctl -u jobsvc -f

[Unit]
Description=Heisenberg jobsvc (FastAPI)
Documentation=https://nimesh08.github.io/quantum-stack/install/server_systemd/
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=jobsvc
Group=jobsvc
WorkingDirectory=/opt/quantum-stack/platform/jobsvc
EnvironmentFile=/etc/quantum-stack/jobsvc.env
ExecStart=/opt/quantum-stack/.venv/bin/uvicorn jobsvc.main:app \
          --host 127.0.0.1 --port 8000 \
          --proxy-headers --forwarded-allow-ips '127.0.0.1'
Restart=on-failure
RestartSec=3

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/jobsvc /var/cache/spinor
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
